# systemd/tunnelgate-renew.service
# tunnelgate-renew.service
[Unit]
Description=Renew TLS certificate
[Service]
Type=oneshot
ExecStart=/usr/local/bin/tunnelgate cert renew
